# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:uvicorn]
command=/app/.venv/bin/uvicorn src.main:app
    --host 127.0.0.1
    --port 8000
    --workers 1
    --timeout-keep-alive 300
directory=/app
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/root",PYTHONPATH="/app:/app/src",PATH="/app/.venv/bin:/root/.local/bin:/usr/local/bin:/usr/bin:/bin"
